# deploy/systemd/onlumis.service
[Unit]
Description=OnLumis Website (Next.js)
After=network.target

[Service]
Type=simple
User=onlumis
Group=onlumis
WorkingDirectory=/var/www/onlumis
EnvironmentFile=/var/www/onlumis/.env.local
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
