# deploy/systemd/wissenskern.service
[Unit]
Description=WissensKern Website (Next.js)
After=network.target

[Service]
Type=simple
User=wissenskern
Group=wissenskern
WorkingDirectory=/var/www/wissenskern
EnvironmentFile=/var/www/wissenskern/.env.local
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
